# yandex-ocr-mcp/deploy/yandex-ocr-mcp.service
[Unit]
Description=yandex-ocr-mcp — MCP-сервер Yandex OCR (Streamable HTTP)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Каталог с исходниками пакета (src/) и установленными зависимостями (node_modules/).
WorkingDirectory=/opt/yandex-ocr-mcp
# Секреты и порт: YANDEX_OCR_API_KEY, MCP_BEARER_TOKEN, PORT.
EnvironmentFile=/opt/yandex-ocr-mcp/.env
ExecStart=/usr/bin/node src/http.ts
Restart=on-failure
RestartSec=3
# Минимальная изоляция (опционально).
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/yandex-ocr-mcp

[Install]
WantedBy=multi-user.target
